# roles/torrent-down/templates/flood.service.j2
[Unit]
Description=Flood - Modern web UI for various torrent clients
After=network-online.target qbittorrent.service
Wants=network-online.target
Requires=qbittorrent.service

[Service]
Type=simple
User=torrent
Group=torrent
WorkingDirectory=/home/torrent
ExecStart=/usr/local/bin/flood --host 0.0.0.0 --port 80 --qburl http://localhost:8080 --qbuser admin --qbpass {{ initial_password }} --allowedpath /mnt/torrents
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target