# roles/navidrome/templates/mnt-music.mount.j2
[Unit]
Description=Mount NAS photos share for Navidrome
After=wait-for-nas.service
Requires=wait-for-nas.service
# This is the key: This mount will not even be attempted until the wait-for-nas service succeeds.
Requires=wait-for-nas.service
After=wait-for-nas.service

[Mount]
What=//{{ nas_ip }}/tank/Media/Music
Where=/mnt/music
Type=cifs
Options=credentials=/etc/navidrome-nas-credentials,uid=1000,gid=1000,iocharset=utf8,vers=3.0,_netdev,x-systemd.device-timeout=30s
